# apps/hyprland/dot-config/systemd/user/wl-wallpaper-update.service
[Unit]
Description=Change wayland wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wl-set-wallpaper
